# additionals/etc/cron.d/elastix.cron
*/5 * * * *  root /usr/bin/php -q /usr/local/elastix/sampler.php
